# deploy/systemd/zephyrus-queue-worker.service
[Unit]
Description=Zephyrus Laravel queue worker
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/Zephyrus
ExecStart=/usr/bin/php artisan queue:work database --queue=integrations,default --sleep=3 --tries=3 --backoff=10,30,60 --timeout=300 --max-time=3600 --memory=256
Restart=always
RestartSec=5
TimeoutStopSec=330
KillSignal=SIGTERM
UMask=0027
PrivateTmp=true
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
